# S1/S1_watchdog_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.3" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>S1 실시간 모니터 감시자 - 로그인/평일 08:00/14:00 실행</Description>
    <URI>\S1_Realtime_Monitor</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-21-940155789-1066660279-1833710528-1001</UserId>
      <LogonType>InteractiveToken</LogonType>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <StartWhenAvailable>true</StartWhenAvailable>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <ExecutionTimeLimit>PT5M</ExecutionTimeLimit>
  </Settings>
  <Triggers>
    <LogonTrigger>
      <Delay>PT1M</Delay>
      <UserId>S-1-5-21-940155789-1066660279-1833710528-1001</UserId>
    </LogonTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-03-12T08:00:00+09:00</StartBoundary>
      <ScheduleByWeek>
        <WeeksInterval>1</WeeksInterval>
        <DaysOfWeek>
          <Monday /><Tuesday /><Wednesday /><Thursday /><Friday />
        </DaysOfWeek>
      </ScheduleByWeek>
    </CalendarTrigger>
    <CalendarTrigger>
      <StartBoundary>2026-03-12T14:00:00+09:00</StartBoundary>
      <ScheduleByWeek>
        <WeeksInterval>1</WeeksInterval>
        <DaysOfWeek>
          <Monday /><Tuesday /><Wednesday /><Thursday /><Friday />
        </DaysOfWeek>
      </ScheduleByWeek>
    </CalendarTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\log\Desktop\Code\S1\S1_monitor_watchdog.bat</Command>
      <WorkingDirectory>C:\Users\log\Desktop\Code\S1</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
